# src/data/coverage.json
{
  "T1566": { "falcon": { "covered": true, "rule": "Phishing attachment / link delivery detection" }, "rapid7": { "covered": true, "rule": "Suspicious email attachment alert" } },
  "T1566.001": { "falcon": { "covered": true, "rule": "Malicious attachment execution from mail client" }, "rapid7": { "covered": false } },
  "T1059": { "falcon": { "covered": true, "rule": "Suspicious command interpreter activity" }, "rapid7": { "covered": false } },
  "T1059.001": { "falcon": { "covered": true, "rule": "Suspicious PowerShell execution" }, "rapid7": { "covered": true, "rule": "Encoded PowerShell command line" } },
  "T1059.003": { "falcon": { "covered": true, "rule": "Suspicious cmd.exe child process chain" }, "rapid7": { "covered": true, "rule": "Living-off-the-land binary execution" } },
  "T1027": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1055": { "falcon": { "covered": true, "rule": "Cross-process memory injection" }, "rapid7": { "covered": false } },
  "T1053": { "falcon": { "covered": true, "rule": "Suspicious scheduled task creation" }, "rapid7": { "covered": true, "rule": "New scheduled task from unusual parent" } },
  "T1053.005": { "falcon": { "covered": true, "rule": "schtasks.exe used to establish persistence" }, "rapid7": { "covered": true, "rule": "Scheduled task persistence detection" } },
  "T1078": { "falcon": { "covered": true, "rule": "Anomalous valid account usage" }, "rapid7": { "covered": true, "rule": "Impossible travel / anomalous login" } },
  "T1078.004": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "Cloud account anomaly detection" } },
  "T1105": { "falcon": { "covered": true, "rule": "Remote file copy / tool transfer detection" }, "rapid7": { "covered": true, "rule": "Ingress tool transfer via network monitoring" } },
  "T1547": { "falcon": { "covered": true, "rule": "Autostart persistence mechanism modified" }, "rapid7": { "covered": false } },
  "T1547.001": { "falcon": { "covered": true, "rule": "Registry Run key modification" }, "rapid7": { "covered": true, "rule": "Startup folder / Run key persistence" } },
  "T1136": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "New local/domain account creation" } },
  "T1098": { "falcon": { "covered": true, "rule": "Suspicious account privilege modification" }, "rapid7": { "covered": true, "rule": "Account manipulation / permission change" } },
  "T1003": { "falcon": { "covered": true, "rule": "Credential dumping tool behavior" }, "rapid7": { "covered": true, "rule": "LSASS access anomaly" } },
  "T1003.001": { "falcon": { "covered": true, "rule": "LSASS memory access from non-standard process" }, "rapid7": { "covered": true, "rule": "Mimikatz-style LSASS dump detection" } },
  "T1021": { "falcon": { "covered": true, "rule": "Lateral movement via remote services" }, "rapid7": { "covered": false } },
  "T1021.001": { "falcon": { "covered": true, "rule": "Suspicious RDP session origin" }, "rapid7": { "covered": true, "rule": "RDP lateral movement detection" } },
  "T1490": { "falcon": { "covered": true, "rule": "Shadow copy / backup deletion" }, "rapid7": { "covered": true, "rule": "vssadmin / wbadmin recovery inhibition" } },
  "T1486": { "falcon": { "covered": true, "rule": "Ransomware-pattern file encryption behavior" }, "rapid7": { "covered": true, "rule": "Mass file modification / entropy spike" } },
  "T1071": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "Anomalous application-layer C2 traffic" } },
  "T1071.001": { "falcon": { "covered": true, "rule": "Suspicious HTTP/HTTPS beaconing" }, "rapid7": { "covered": true, "rule": "Web protocol C2 beacon detection" } },
  "T1041": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1048": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "Exfiltration over DNS / alternate protocol" } },
  "T1567": { "falcon": { "covered": true, "rule": "Data transfer to cloud storage service" }, "rapid7": { "covered": false } },
  "T1082": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1057": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1018": { "falcon": { "covered": true, "rule": "Network scanning / host discovery behavior" }, "rapid7": { "covered": true, "rule": "Internal network sweep detection" } },
  "T1016": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1049": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1087": { "falcon": { "covered": true, "rule": "Account enumeration via net.exe / PowerShell" }, "rapid7": { "covered": false } },
  "T1069": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1560": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "Archive utility staging collected data" } },
  "T1140": { "falcon": { "covered": true, "rule": "Script deobfuscation / decode behavior" }, "rapid7": { "covered": false } },
  "T1204": { "falcon": { "covered": true, "rule": "User-initiated malicious execution" }, "rapid7": { "covered": true, "rule": "Suspicious user-executed file alert" } },
  "T1204.002": { "falcon": { "covered": true, "rule": "Malicious file opened by user" }, "rapid7": { "covered": true, "rule": "Malicious document execution chain" } },
  "T1110": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "Brute-force authentication attempts" } },
  "T1110.003": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "Password spraying across multiple accounts" } },
  "T1112": { "falcon": { "covered": true, "rule": "Suspicious registry modification" }, "rapid7": { "covered": false } },
  "T1070": { "falcon": { "covered": true, "rule": "Indicator/log removal behavior" }, "rapid7": { "covered": true, "rule": "Event log clearing detection" } },
  "T1070.004": { "falcon": { "covered": true, "rule": "Suspicious file deletion post-execution" }, "rapid7": { "covered": false } },
  "T1036": { "falcon": { "covered": true, "rule": "Masquerading as legitimate system binary" }, "rapid7": { "covered": true, "rule": "Process name / path mismatch detection" } },
  "T1036.005": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1543": { "falcon": { "covered": true, "rule": "New system service creation for persistence" }, "rapid7": { "covered": false } },
  "T1543.003": { "falcon": { "covered": true, "rule": "Suspicious Windows service installation" }, "rapid7": { "covered": true, "rule": "Service creation persistence detection" } },
  "T1569": { "falcon": { "covered": false }, "rapid7": { "covered": false } },
  "T1569.002": { "falcon": { "covered": true, "rule": "sc.exe / PsExec service execution" }, "rapid7": { "covered": true, "rule": "Remote service execution detection" } },
  "T1219": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "Remote access tool (AnyDesk/TeamViewer) usage" } },
  "T1133": { "falcon": { "covered": false }, "rapid7": { "covered": true, "rule": "External remote service access anomaly" } },
  "T1190": { "falcon": { "covered": true, "rule": "Exploitation attempt against public-facing app" }, "rapid7": { "covered": true, "rule": "Web application exploit signature" } }
}
